# rapid/deploy/rapid.service
# deploy/rapid.service — systemd unit file for RAPID
#
# Install:
#   sudo cp deploy/rapid.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now rapid
#
# Logs: journalctl -u rapid -f

[Unit]
Description=RAPID — RAG Application for Private Instant Deployment
After=network.target

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/opt/rapid/rapid
EnvironmentFile=/opt/rapid/rapid/.env

# Activate virtualenv and start Gunicorn with UvicornWorker
ExecStart=/opt/rapid/venv/bin/gunicorn \
    -c deploy/gunicorn.conf.py \
    -k uvicorn.workers.UvicornWorker \
    main:app

ExecReload=/bin/kill -s HUP $MAINPID

# Graceful shutdown: finish in-flight requests
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=35

# Restart policy
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/rapid /tmp

[Install]
WantedBy=multi-user.target
